# build/nginx/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/app/run/supervisord.pid

[unix_http_server]
file=/app/run/supervisord.sock
username = dummy
password = dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///app/run/supervisord.sock
username = dummy
password = dummy

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
autostart=false
startretries=0
startsecs=3
priority=20

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -e /app/run/nginx.error.log -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
autostart=false
startretries=0
startsecs=3
priority=30

[program:scheduler]
command=php artisan schedule:work
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autorestart=true
autostart=false
startsecs=3
priority=40

[program:entrypoint]
command=/app/Backend/entrypoint.sh
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autorestart=false
autostart=true
startsecs=0
startretries=1
priority=10
